# deploy/supervisor/conf.d/litellm.conf
[program:litellm]
command=%(ENV_HOME)s/.local/bin/litellm --config %(ENV_HOME)s/.litellm/config.yaml --port 4000
directory=%(ENV_HOME)s
autostart=false
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/litellm/logs/litellm.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/litellm/logs/litellm_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
