# deploy/supervisor/conf.d/lyra_tts.conf
[program:lyra_tts]
command=%(ENV_HOME)s/projects/lyra/supervisor/scripts/run_adapter.sh tts
directory=%(ENV_HOME)s/projects/lyra
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/projects/lyra/.venv/bin:%(ENV_PATH)s"
autostart=false
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=75
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/lyra/logs/lyra_tts.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/lyra/logs/lyra_tts_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
