# deploy/supervisor/conf.d/lyra_imagecli_gen.conf
[program:lyra_imagecli_gen]
command=imagecli nats-serve
directory=%(ENV_HOME)s/projects/lyra
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/projects/lyra/.venv/bin:%(ENV_PATH)s",NATS_NKEY_SEED_PATH="%(ENV_HOME)s/.lyra/nkeys/image-worker.seed",NATS_URL="nats://127.0.0.1:4222"
priority=200
autostart=false
autorestart=true
startsecs=20
startretries=3
stopwaitsecs=75
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/lyra/logs/lyra_imagecli_gen.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/lyra/logs/lyra_imagecli_gen_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
