# deploy/supervisor/conf.d/lyra_hub.conf
[program:lyra_hub]
command=%(ENV_HOME)s/projects/lyra/deploy/supervisor/scripts/run_hub.sh
directory=%(ENV_HOME)s/projects/lyra
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/projects/lyra/.venv/bin:%(ENV_PATH)s",LYRA_HEALTH_PORT="8443",NATS_NKEY_SEED_PATH="%(ENV_HOME)s/.lyra/nkeys/hub.seed"
autostart=false
priority=100
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=75
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/lyra/logs/lyra_hub.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/lyra/logs/lyra_hub_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
