# deploy/supervisor/conf.d/voicecli_stt.conf
[program:voicecli_stt]
command=%(ENV_HOME)s/projects/lyra/deploy/supervisor/scripts/run_voicecli.sh stt
directory=%(ENV_HOME)s/projects/voiceCLI
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/projects/voiceCLI/.venv/bin:/usr/local/bin:/usr/bin:/bin",NATS_URL="nats://127.0.0.1:4222",NATS_NKEY_SEED_PATH="%(ENV_HOME)s/.lyra/nkeys/voice-stt.seed",LYRA_STT_ENGINE="whisper"
autostart=false
autorestart=true
# startsecs: supervisor marks RUNNING after this many seconds of uptime (time-based,
# not heartbeat-based). Tune only if the process is crash-looping during GPU warm-up.
# Re-measure first-ready time on Machine 1 first boot and document in PR/runbook (T15).
startsecs=15
startretries=3
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/voicecli/logs/voicecli_stt.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/voicecli/logs/voicecli_stt_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
