# deploy/supervisor/conf.d/voicecli_tts.conf
[program:voicecli_tts]
command=%(ENV_HOME)s/projects/lyra/deploy/supervisor/scripts/run_voicecli.sh tts
directory=%(ENV_HOME)s/projects/voiceCLI
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/projects/voiceCLI/.venv/bin:/usr/local/bin:/usr/bin:/bin",NATS_URL="nats://127.0.0.1:4222",NATS_NKEY_SEED_PATH="%(ENV_HOME)s/.lyra/nkeys/voice-tts.seed",LYRA_TTS_ENGINE="chatterbox"
autostart=false
autorestart=true
# startsecs: see voicecli_stt.conf for semantics. TTS Chatterbox model load may
# need this raised to 15 on first boot — re-tune per T15 measurement.
startsecs=10
startretries=3
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/voicecli/logs/voicecli_tts.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/voicecli/logs/voicecli_tts_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
