# supervisor/conf.d/lyra-monitoring.conf
[program:lyra_monitoring]
; Safety-net health check — runs every 4 hours (14400s).
; Real-time monitoring is handled by EventAggregator (issue #44).
; This process catches anything the event system misses (e.g. hub is down).
;
; Loop: run checks → sleep 4h → repeat. Never exits, so supervisord keeps it alive.
command=/bin/sh -c 'while true; do %(ENV_HOME)s/projects/lyra/.venv/bin/python -m lyra.monitoring; sleep 14400; done'
directory=%(ENV_HOME)s/projects/lyra
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_HOME)s/projects/lyra/.venv/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
startsecs=2
startretries=3
stopwaitsecs=15
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/lyra/logs/lyra_monitoring.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=%(ENV_HOME)s/.local/state/lyra/logs/lyra_monitoring_error.log
stderr_logfile_maxbytes=2MB
stderr_logfile_backups=2
