# deploy/lyra.service
[Unit]
Description=Lyra supervisord (hub, telegram, discord)
After=network-online.target nats.service
Wants=network-online.target
Requires=nats.service

[Service]
Type=forking
PIDFile=%h/projects/lyra/deploy/supervisor/supervisord.pid
ExecStart=%h/projects/lyra/deploy/supervisor/start.sh --all
ExecStop=%h/projects/lyra/deploy/supervisor/supervisorctl.sh shutdown
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
