# deploy/lyra-monitor.timer
[Unit]
Description=Run Lyra health monitor periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec={{MONITOR_INTERVAL:-5min}}
Persistent=true

[Install]
WantedBy=timers.target
