# deploy/lyra.service
[Unit]
Description=Lyra supervisord (hub, telegram, discord)
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=%h/projects/supervisord.pid
ExecStart=%h/projects/scripts/start.sh --all
ExecStop=%h/projects/scripts/supervisorctl.sh shutdown
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
