# deploy/lyra.service
[Unit]
Description=Lyra supervisord (hub, telegram, discord)
After=graphical-session.target

[Service]
Type=forking
PIDFile=%h/projects/lyra/deploy/supervisor/supervisord.pid
ExecStart=%h/projects/lyra/deploy/supervisor/start.sh
ExecStop=%h/.local/bin/supervisorctl -c %h/projects/lyra/deploy/supervisor/supervisord.conf shutdown
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
