# deploy/lyra-monitor.service
[Unit]
Description=Lyra health monitor (Layer 1 pre-check + Layer 2 LLM escalation)
After=network-online.target lyra.service
Wants=network-online.target

[Service]
Type=oneshot
User={{LYRA_USER:-mickael}}
WorkingDirectory={{LYRA_PROJECT_DIR}}
EnvironmentFile=%h/.lyra/.env
ExecStart={{LYRA_PROJECT_DIR}}/.venv/bin/python -m lyra.monitoring
TimeoutStartSec=60

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=%h/.lyra
