# deploy/lyra.service
[Unit]
Description=Start machine-level supervisord (Lyra tenant)
After=network-online.target nats.service
Wants=network-online.target nats.service

[Service]
Type=forking
PIDFile=%h/projects/supervisord.pid
ExecStart=%h/projects/scripts/start.sh --all
ExecStop=%h/projects/scripts/supervisorctl.sh shutdown
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
