# deploy/quadlet/discord.container
[Unit]
Description=Lyra Discord adapter
After=lyra-hub.service
Requires=lyra-hub.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
Image=localhost/lyra:latest
ContainerName=lyra-discord
Network=lyra.network
EnvironmentFile=%h/projects/lyra/.env.discord
Environment=NATS_URL=nats://lyra-nats:4222
Volume=lyra-nkey-discord-adapter.volume:/run/secrets/discord-adapter.seed:ro,z
Environment=NATS_NKEY_SEED_PATH=/run/secrets/discord-adapter.seed
Volume=lyra-data.volume:/home/lyra/.lyra:ro,z
Volume=lyra-config.volume:/app/config.toml:ro,z
Exec=lyra adapter discord

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
