# deploy/quadlet/hub.container
[Unit]
Description=Lyra Hub
After=nats.service
Requires=nats.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
Image=localhost/lyra:latest
ContainerName=lyra-hub
Network=lyra.network
EnvironmentFile=%h/projects/lyra/.env.hub
# .env.hub contains only hub-scoped vars (ANTHROPIC_API_KEY, LYRA_HEALTH_SECRET,
# LYRA_LOG_DIR, ROXABI_VAULT_DIR). Adapter tokens (Telegram, Discord, TTS/STT)
# are excluded — hub does not need them. See deploy/quadlet/.env.hub.example.
# NATS_URL is overridden inline below; cutover = update .env.hub directly.
Environment=NATS_URL=nats://lyra-nats:4222
Volume=lyra-nkey-hub.volume:/run/secrets/hub.seed:ro,z
Environment=NATS_NKEY_SEED_PATH=/run/secrets/hub.seed
Volume=lyra-data.volume:/home/lyra/.lyra:z
Volume=lyra-config.volume:/app/config.toml:ro,z
# Health endpoint — probed by lyra-monitor via host loopback (localhost:8443).
PublishPort=127.0.0.1:8443:8443
Exec=lyra hub

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
