# deploy/quadlet/lyra-telegram.container
[Unit]
Description=Lyra Telegram adapter
# After= orders startup; no Requires=/BindsTo= so a hub restart does not
# propagate a stop to the adapter. Adapters reconnect to NATS on their own
# once the hub is back; tearing them down blocks auto-recovery.
After=lyra-hub.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
# Image digest pinning for localhost/lyra:latest tracks on a follow-up issue
# (CI tag-by-SHA strategy).
Image=ghcr.io/roxabi/lyra:staging
Label=io.containers.autoupdate=registry
ContainerName=lyra-telegram
Network=roxabi.network
# Hardening — defense in depth (see issue #652).
NoNewPrivileges=true
ReadOnly=true
DropCapability=all
# ADR-054 Decision 2: UserNS=keep-id maps container UID to host UID 1000 (mickael),
# which owns ~/.lyra/*. Replaces User=lyra (UID 1500) which could not read host files.
UserNS=keep-id:uid=1500,gid=1500
# ADR-054 Decision 4: bot tokens come from config.db via LyraKeyring (since #417).
# No EnvironmentFile needed for tokens — telegram.env.example deleted.
Environment=NATS_URL=nats://lyra-nats:4222
# ADR-054 Decision 5 (revised): nkey seed delivered as Podman secret.
Secret=lyra-nkey-telegram-adapter,type=mount,target=telegram-adapter.seed,mode=0400,uid=1500,gid=1500
Environment=NATS_NKEY_SEED_PATH=/run/secrets/telegram-adapter.seed
# ADR-054 Decision 3: rw (adapters write discord.db, message_index.db, auth.db pairings).
Volume=lyra-data.volume:/home/lyra/.lyra:z
Volume=lyra-logs.volume:/home/lyra/.local/state/lyra/logs:z
Environment=LYRA_LOG_DIR=/home/lyra/.local/state/lyra/logs
# Inline single-file bind mount — see rationale in lyra-hub.container.
Volume=%h/.lyra/config.toml:/app/config.toml:ro,z
Exec=lyra adapter telegram

[Service]
Restart=on-failure
RestartSec=5
# Cap at 2 full cores — prevents a runaway worker from starving hub + NATS.
CPUQuota=200%

[Install]
WantedBy=default.target
