# deploy/quadlet/lyra-discord.container
[Unit]
Description=Lyra Discord adapter
# After= orders startup; no Requires=/BindsTo= so a hub restart does not
# propagate a stop to the adapter. Adapters reconnect to NATS on their own
# once the hub is back; tearing them down blocks auto-recovery.
After=lyra-hub.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
# Image digest pinning for localhost/lyra:latest tracks on a follow-up issue
# (CI tag-by-SHA strategy).
Image=ghcr.io/roxabi/lyra:staging
Label=io.containers.autoupdate=registry
ContainerName=lyra-discord
Network=roxabi.network
# Hardening — defense in depth (see issue #652).
NoNewPrivileges=true
ReadOnly=true
DropCapability=all
# ADR-054 Decision 2: UserNS=keep-id maps container UID to host UID 1000 (mickael),
# which owns ~/.lyra/*. Replaces User=lyra (UID 1500) which could not read host files.
UserNS=keep-id:uid=1500,gid=1500
# ADR-054 Decision 4: bot tokens come from config.db via LyraKeyring (since #417).
# No EnvironmentFile needed for tokens — discord.env.example deleted.
Environment=NATS_URL=nats://lyra-nats:4222
# ADR-054 Decision 5 (revised): nkey seed delivered as Podman secret.
Secret=lyra-nkey-discord-adapter,type=mount,target=discord-adapter.seed,mode=0400,uid=1500,gid=1500
Environment=NATS_NKEY_SEED_PATH=/run/secrets/discord-adapter.seed
# Per-file bind mounts — adapters only access keyring.key, config.db, turns.db, and discord.db.
# Full lyra-data.volume removed (#943) — blast radius reduction on adapter compromise.
# Inline binds used for single files (named .volume units fail for files; see hub.container note).
# :z (shared label) not :Z — hub's directory-level lyra-data.volume relabels ~/.lyra/,
# making :Z (private) incompatible. Shared label is intentional (#943).
Volume=%h/.lyra/keyring.key:/home/lyra/.lyra/keyring.key:ro,z
Volume=%h/.lyra/config.db:/home/lyra/.lyra/config.db:ro,z
Volume=%h/.lyra/turns.db:/home/lyra/.lyra/turns.db:ro,z
Volume=%h/.lyra/discord.db:/home/lyra/.lyra/discord.db:z
Volume=lyra-logs.volume:/home/lyra/.local/state/lyra/logs:z
Environment=LYRA_LOG_DIR=/home/lyra/.local/state/lyra/logs
# Inline single-file bind mount — see rationale in lyra-hub.container.
Volume=%h/.lyra/config.toml:/app/config.toml:ro,z
Exec=lyra adapter discord

[Service]
Restart=on-failure
RestartSec=5
# Cap at 2 full cores — prevents a runaway worker from starving hub + NATS.
CPUQuota=200%
# Ensure SQLite files exist before Podman binds them — missing source path causes
# Podman to auto-create a directory, breaking SQLite open on first boot.
ExecStartPre=/bin/bash -c 'test -e %h/.lyra/turns.db || touch %h/.lyra/turns.db'
ExecStartPre=/bin/bash -c 'test -e %h/.lyra/discord.db || touch %h/.lyra/discord.db'

[Install]
WantedBy=default.target
