# deploy/quadlet/lyra-hub.container
[Unit]
Description=Lyra Hub
After=nats.service
Requires=nats.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
# Image digest pinning for localhost/lyra:latest tracks on a follow-up issue
# (CI tag-by-SHA strategy). Upstream registry images (e.g. nats.container)
# are pinned by digest directly.
Image=ghcr.io/roxabi/lyra:staging
ContainerName=lyra-hub
Network=roxabi.network
# Hardening — defense in depth (see issue #652).
NoNewPrivileges=true
ReadOnly=true
DropCapability=all
# ADR-054 Decision 2: UserNS=keep-id maps container UID to host UID 1000 (mickael),
# which owns ~/.lyra/*. Replaces User=lyra (UID 1500) which could not read host files.
UserNS=keep-id:uid=1500,gid=1500
EnvironmentFile=%h/.lyra/env/hub.env
# hub.env contains only hub-scoped vars (LYRA_HEALTH_SECRET,
# LYRA_LOG_DIR, ROXABI_VAULT_DIR). Adapter tokens (Telegram, Discord, TTS/STT)
# are excluded — hub does not need them. See deploy/quadlet/hub.env.example.
# NATS_URL is overridden inline below; cutover = update .env.hub directly.
Environment=NATS_URL=nats://lyra-nats:4222
# ADR-054 Decision 5 (revised): nkey seed delivered as Podman secret.
Secret=lyra-nkey-hub,type=mount,target=hub.seed,mode=0400,uid=1500,gid=1500
Environment=NATS_NKEY_SEED_PATH=/run/secrets/hub.seed
Volume=lyra-data.volume:/home/lyra/.lyra:z
Volume=lyra-logs.volume:/home/lyra/.local/state/lyra/logs:z
Environment=LYRA_LOG_DIR=/home/lyra/.local/state/lyra/logs
# Inline single-file bind mount — Quadlet .volume units wrapping a single file
# either fail with loop-device errors or are silently downgraded to empty named
# volumes (IsADirectoryError at read time). Inline Volume= does a direct bind.
Volume=%h/.lyra/config.toml:/app/config.toml:ro,z
# claude CLI — selective mounts (not the whole ~/.claude dir).
# Only config, credentials, memory, and skills are shared — ephemeral state
# (sessions, tasks, teams, history, telemetry, cache) stays host-only.
# UserNS=keep-id maps UID 1500 → host UID 1000 (mickael).
#   rw  — CLI writes tokens + project memory at runtime
#   ro  — config/skills read at startup, no container writes
#   z   — shared SELinux label so host CLI retains access
Volume=%h/.claude/.credentials.json:/home/lyra/.claude/.credentials.json:z
Volume=%h/.claude/settings.json:/home/lyra/.claude/settings.json:ro,z
Volume=%h/.claude/settings.local.json:/home/lyra/.claude/settings.local.json:ro,z
Volume=%h/.claude/CLAUDE.md:/home/lyra/.claude/CLAUDE.md:ro,z
Volume=%h/.claude/projects:/home/lyra/.claude/projects:z
Volume=%h/.claude/plugins:/home/lyra/.claude/plugins:ro,z
Volume=%h/.claude/skills:/home/lyra/.claude/skills:ro,z
Volume=%h/.claude/shared:/home/lyra/.claude/shared:ro,z
Volume=%h/.claude/.git:/home/lyra/.claude/.git:ro,z
Volume=%h/.claude.json:/home/lyra/.claude.json:z
# Workspace — mirrors host ~/projects so claude subprocess picks up
# ~/projects/CLAUDE.md (project index) and per-project instructions.
# This is the default cwd for agent claude processes (config.toml [defaults]).
Volume=%h/projects:/home/lyra/projects:z
# Default cwd for claude subprocesses — ~/projects (not /app).
Environment=LYRA_CLAUDE_CWD=/home/lyra/projects
# Health endpoint — probed by lyra-monitor via host loopback (localhost:8443).
# Bind to 0.0.0.0 inside container so PublishPort forwarding reaches uvicorn;
# host exposure is restricted to 127.0.0.1 by PublishPort below.
Environment=LYRA_HEALTH_HOST=0.0.0.0
PublishPort=127.0.0.1:8443:8443
Exec=lyra hub

[Service]
Restart=on-failure
RestartSec=5
# Cap at 2 full cores — prevents a runaway worker from starving hub + NATS.
CPUQuota=200%
# Pre-create single-file bind-mount sources if absent — prevents Podman from
# auto-creating them as directories (IsADirectoryError at read time).
ExecStartPre=/bin/bash -c 'test -e %h/.claude.json || echo "{}" > %h/.claude.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/.credentials.json || echo "{}" > %h/.claude/.credentials.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/settings.json || echo "{}" > %h/.claude/settings.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/settings.local.json || echo "{}" > %h/.claude/settings.local.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/CLAUDE.md || touch %h/.claude/CLAUDE.md'

[Install]
WantedBy=default.target
