# deploy/quadlet/lyra-clipool.container
[Unit]
Description=Lyra CliPool NATS worker
After=lyra-nats.service
Requires=lyra-nats.service

[Container]
Image=ghcr.io/roxabi/lyra:staging
Label=io.containers.autoupdate=registry
ContainerName=lyra-clipool
Network=roxabi.network
NoNewPrivileges=true
ReadOnly=true
DropCapability=all
UserNS=keep-id:uid=1500,gid=1500
Environment=NATS_URL=nats://lyra-nats:4222
Secret=lyra-nkey-clipool,type=mount,target=clipool.seed,mode=0400,uid=1500,gid=1500
Environment=NATS_NKEY_SEED_PATH=/run/secrets/clipool.seed
Volume=lyra-logs.volume:/home/lyra/.local/state/lyra/logs:z
Environment=LYRA_LOG_DIR=/home/lyra/.local/state/lyra/logs
Volume=%h/.lyra/config.toml:/app/config.toml:ro,z
Volume=%h/.claude/.credentials.json:/home/lyra/.claude/.credentials.json:z
Volume=%h/.claude/settings.json:/home/lyra/.claude/settings.json:ro,z
Volume=%h/.claude/settings.local.json:/home/lyra/.claude/settings.local.json:ro,z
Volume=%h/.claude/CLAUDE.md:/home/lyra/.claude/CLAUDE.md:ro,z
Volume=%h/.claude/projects:/home/lyra/.claude/projects:z
Volume=%h/.claude/plugins:/home/lyra/.claude/plugins:ro,z
Volume=%h/.claude/skills:/home/lyra/.claude/skills:ro,z
Volume=%h/.claude/shared:/home/lyra/.claude/shared:ro,z
Volume=%h/.claude/.git:/home/lyra/.claude/.git:ro,z
Volume=%h/.claude.json:/home/lyra/.claude.json:z
Volume=%h/projects:/home/lyra/projects:z
Environment=LYRA_CLAUDE_CWD=/home/lyra/projects
ExecStartPre=/bin/bash -c 'test -e %h/.claude.json || echo "{}" > %h/.claude.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/.credentials.json || echo "{}" > %h/.claude/.credentials.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/settings.json || echo "{}" > %h/.claude/settings.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/settings.local.json || echo "{}" > %h/.claude/settings.local.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/CLAUDE.md || touch %h/.claude/CLAUDE.md'
Exec=lyra clipool

[Service]
Restart=on-failure
RestartSec=5
CPUQuota=400%

[Install]
WantedBy=default.target
