# deploy/quadlet/lyra-clipool.container
[Unit]
Description=Lyra CliPool Worker
After=lyra-nats.service
# No Requires= — clipool connects directly to NATS; hub restarts must not cascade

[Container]
Image=ghcr.io/roxabi/lyra:staging
Label=io.containers.autoupdate=registry
ContainerName=lyra-clipool
Network=roxabi.network
NoNewPrivileges=true
ReadOnly=true
DropCapability=all
UserNS=keep-id:uid=1500,gid=1500
Environment=NATS_URL=nats://lyra-nats:4222
Volume=lyra-logs.volume:/home/lyra/.local/state/lyra/logs:z
Environment=LYRA_LOG_DIR=/home/lyra/.local/state/lyra/logs
Secret=lyra-nkey-clipool-worker,type=mount,target=clipool-worker.seed,mode=0400,uid=1500,gid=1500
Environment=NATS_NKEY_SEED_PATH=/run/secrets/clipool-worker.seed
# ~/.claude/ — directory mount (not file mounts) so Syncthing's atomic renames
# (unlink old inode → rename new) remain transparent to the container.
# Single-file bind-mounts bind the inode at mount time; after an atomic replace
# the container holds a Links:0 ghost and never sees the refreshed token.
Volume=%h/.claude:/home/lyra/.claude:z
Volume=%h/.claude.json:/home/lyra/.claude.json:z
# Workspace — mirrors host ~/projects so claude subprocess picks up per-project CLAUDE.md
Volume=%h/projects:/home/lyra/projects:z
Environment=LYRA_CLAUDE_CWD=/home/lyra/projects
Exec=lyra adapter clipool

[Service]
Restart=on-failure
RestartSec=5
TimeoutStopSec=120
CPUQuota=400%
ExecStartPre=/bin/bash -c 'test -e %h/.claude.json || echo "{}" > %h/.claude.json'
ExecStartPre=/bin/bash -c 'test -e %h/.claude/.credentials.json || echo "{}" > %h/.claude/.credentials.json'

[Install]
WantedBy=default.target
