# deploy/quadlet/lyra-hub.container
[Unit]
Description=Lyra Hub
After=nats.service
Requires=nats.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
# Image digest pinning for localhost/lyra:latest tracks on a follow-up issue
# (CI tag-by-SHA strategy). Upstream registry images (e.g. nats.container)
# are pinned by digest directly.
Image=localhost/lyra:latest
ContainerName=lyra-hub
Network=lyra.network
# Hardening — defense in depth (see issue #652).
NoNewPrivileges=true
ReadOnly=true
DropCapability=all
# ADR-054 Decision 2: UserNS=keep-id maps container UID to host UID 1000 (mickael),
# which owns ~/.lyra/*. Replaces User=lyra (UID 1500) which could not read host files.
UserNS=keep-id:uid=1500,gid=1500
EnvironmentFile=%h/.lyra/env/hub.env
# hub.env contains only hub-scoped vars (LYRA_HEALTH_SECRET,
# LYRA_LOG_DIR, ROXABI_VAULT_DIR). Adapter tokens (Telegram, Discord, TTS/STT)
# are excluded — hub does not need them. See deploy/quadlet/hub.env.example.
# NATS_URL is overridden inline below; cutover = update .env.hub directly.
Environment=NATS_URL=nats://lyra-nats:4222
# ADR-054 Decision 5 (revised): nkey seed delivered as Podman secret.
Secret=lyra-nkey-hub,type=mount,target=hub.seed,mode=0400,uid=1500,gid=1500
Environment=NATS_NKEY_SEED_PATH=/run/secrets/hub.seed
Volume=lyra-data.volume:/home/lyra/.lyra:z
Volume=lyra-logs.volume:/home/lyra/.local/state/lyra/logs:z
Environment=LYRA_LOG_DIR=/home/lyra/.local/state/lyra/logs
# Inline single-file bind mount — Quadlet .volume units wrapping a single file
# either fail with loop-device errors or are silently downgraded to empty named
# volumes (IsADirectoryError at read time). Inline Volume= does a direct bind.
Volume=%h/.lyra/config.toml:/app/config.toml:ro,z
# claude CLI credentials — host ~/.claude is bind-mounted read-write so the
# CLI can refresh tokens. UserNS=keep-id maps UID 1500 → host UID 1000 (mickael).
# ~/.claude.json is a separate top-level config file the CLI looks up outside
# of ~/.claude/ — must be bind-mounted inline (same pattern as config.toml).
# Both mounts are :rw,z by design — host ↔ container share full state
# (sessions, projects, skills, tokens) so a refresh on either side propagates
# immediately, mirror-like. Scalpel split (credential-only rw, history ro)
# was considered and rejected: the whole-dir share is the intended behavior.
#   rw  — CLI writes tokens + project/session state at runtime
#   z   — shared SELinux label so host CLI retains access; :Z would isolate
#         per-container and break host CLI on SELinux-enforcing systems
Volume=%h/.claude:/home/lyra/.claude:z
Volume=%h/.claude.json:/home/lyra/.claude.json:z
# Health endpoint — probed by lyra-monitor via host loopback (localhost:8443).
# Bind to 0.0.0.0 inside container so PublishPort forwarding reaches uvicorn;
# host exposure is restricted to 127.0.0.1 by PublishPort below.
Environment=LYRA_HEALTH_HOST=0.0.0.0
PublishPort=127.0.0.1:8443:8443
Exec=lyra hub

[Service]
Restart=on-failure
RestartSec=5
# Cap at 2 full cores — prevents a runaway worker from starving hub + NATS.
CPUQuota=200%
# Pre-create ~/.claude.json as empty JSON if absent — prevents Podman from
# auto-creating the bind-mount source as a directory (IsADirectoryError at
# CLI read time). Touch alone leaves a 0-byte file that json parsers reject;
# '{}' is the minimal valid document the CLI can read and then populate.
ExecStartPre=/bin/bash -c 'test -e %h/.claude.json || echo "{}" > %h/.claude.json'

[Install]
WantedBy=default.target
