# deploy/quadlet/lyra-discord.container
[Unit]
Description=Lyra Discord adapter
# After= orders startup; no Requires=/BindsTo= so a hub restart does not
# propagate a stop to the adapter. Adapters reconnect to NATS on their own
# once the hub is back; tearing them down blocks auto-recovery.
After=lyra-hub.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
# Image digest pinning for localhost/lyra:latest tracks on a follow-up issue
# (CI tag-by-SHA strategy).
Image=localhost/lyra:latest
ContainerName=lyra-discord
Network=lyra.network
# Hardening — defense in depth (see issue #652).
NoNewPrivileges=true
ReadOnly=true
DropCapability=all
User=lyra
# Scoped env file — copy from discord.env.example, fill in bot token only.
# NATS_URL and NATS_NKEY_SEED_PATH are set inline below.
EnvironmentFile=%h/.lyra/env/discord.env
Environment=NATS_URL=nats://lyra-nats:4222
Volume=lyra-nkey-discord-adapter.volume:/run/secrets/discord-adapter.seed:ro,z
Environment=NATS_NKEY_SEED_PATH=/run/secrets/discord-adapter.seed
Volume=lyra-data.volume:/home/lyra/.lyra:ro,z
Volume=lyra-logs.volume:/home/lyra/.local/state/lyra/logs:z
Environment=LYRA_LOG_DIR=/home/lyra/.local/state/lyra/logs
Volume=lyra-config.volume:/app/config.toml:ro,z
Exec=lyra adapter discord

[Service]
Restart=on-failure
RestartSec=5
# Cap at 2 full cores — prevents a runaway worker from starving hub + NATS.
CPUQuota=200%

[Install]
WantedBy=default.target
