# deploy/quadlet/telegram.container
[Unit]
Description=Lyra Telegram adapter
After=lyra-hub.service
Requires=lyra-hub.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
Image=localhost/lyra:latest
ContainerName=lyra-telegram
Network=lyra.network
EnvironmentFile=%h/projects/lyra/.env.telegram
Environment=NATS_URL=nats://lyra-nats:4222
Volume=lyra-nkey-telegram-adapter.volume:/run/secrets/telegram-adapter.seed:ro,z
Environment=NATS_NKEY_SEED_PATH=/run/secrets/telegram-adapter.seed
Volume=lyra-data.volume:/home/lyra/.lyra:ro,z
Volume=lyra-config.volume:/app/config.toml:ro,z
Exec=lyra adapter telegram

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
